# meta-propertycore/recipes-propertycore/propertycore-engine/files/propertycore-engine.service
[Unit]
Description=PropertyCore Automation Engine
Documentation=https://github.com/donzeg/PropertyCore
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
ExecStart=/usr/bin/propertycore-engine
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=propertycore-engine

# Persistence — systemd creates /var/lib/propertycore/ before starting the service
StateDirectory=propertycore
StateDirectoryMode=0755

# Engine configuration — override per deployment via /etc/propertycore/engine.conf
Environment=MQTT_BROKER=localhost:1883

[Install]
WantedBy=multi-user.target
